# com.linear-claude.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.linear-claude.automation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/rapidinnovation/Desktop/linear-automation-python/run_once.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rapidinnovation/Desktop/linear-automation-python</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/rapidinnovation/Desktop/linear-automation-python/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rapidinnovation/Desktop/linear-automation-python/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/rapidinnovation</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
